# credential-sniffer/sniffer.service
[Unit]
Description=Credentials Sniffer service

[Service]
ExecStart=/home/umbox/credential-sniffer/sniffer.sh

[Install]
WantedBy=multi-user.target